# monitor.service
[Unit]
Description=Server Monitor - Dashboard de Monitoramento
After=network.target docker.service
Wants=docker.service

# NOTE: @@USER@@ and @@REMOTE_DIR@@ are placeholders automatically filled in
# by deploy.py at deploy time. If installing manually, replace them yourself
# (e.g. USER=gabe, REMOTE_DIR=/home/gabe/server-monitor) before copying this
# file to /etc/systemd/system/monitor.service.

[Service]
Type=simple
User=@@USER@@
WorkingDirectory=@@REMOTE_DIR@@
ExecStart=@@REMOTE_DIR@@/venv/bin/python app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

# Limites de recursos: o monitor roda direto no host (fora do Docker) e não
# deve competir por recursos com os serviços que ele monitora.
MemoryMax=150M
MemoryHigh=100M
CPUQuota=25%
Nice=10
TasksMax=64

[Install]
WantedBy=multi-user.target
